# deploy/systemd/ak-proxy.service
[Unit]
Description=AK Proxy Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/ak-proxy

# Load environment variables from file (optional).
# Put NOTIFY_CENTER_INTERNAL_SECRET in /etc/ak-proxy.env
EnvironmentFile=-/etc/ak-proxy.env

# Generate missing deployment secrets, then run the server with the generated values.
# Existing values are preserved; an explicitly empty signing key fails closed.
ExecStart=/home/ubuntu/ak-proxy/venv/bin/python /home/ubuntu/ak-proxy/deploy/systemd/ak_proxy_bootstrap.py
Restart=always
RestartSec=2

# Hardening (keep minimal to avoid breaking existing behavior)
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
